# whispermodel.service
[Unit]
Description=whispermodel — system-wide Whisper STT daemon (whisper.cpp)
After=graphical-session.target

[Service]
Type=simple
# Assumes the repo lives at ~/programming/whispermodel. If you cloned it elsewhere,
# edit this path (or symlink the repo there).
ExecStart=%h/programming/whispermodel/whispermodel.sh
Restart=on-failure
RestartSec=2
# The model stays loaded for the whole session so transcription is instant. For idle-unload
# instead, switch to socket activation / add an idle timeout in the launcher.

[Install]
WantedBy=default.target
